# backend/mjg-backend.service
[Unit]
Description=麻将馆记账系统后端服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/data/mjg/backend
ExecStart=/data/mjg/backend/start.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target








